# ops/mac/launchd/ai.tmc.public_research_ledger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>ai.tmc.public_research_ledger</string>
  <key>ProgramArguments</key><array>
    <string>/Users/mikicourage/.hermes/hermes-agent/venv/bin/python</string>
    <string>/Users/mikicourage/TradingMapClaw/ops_migration/public_outbox.py</string>
    <string>verification-ledger</string><string>--prepare</string>
  </array>
  <key>StartCalendarInterval</key><dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key><string>/Users/mikicourage/TradingMapClaw/ops_migration/logs/public_research_ledger.out.log</string>
  <key>StandardErrorPath</key><string>/Users/mikicourage/TradingMapClaw/ops_migration/logs/public_research_ledger.err.log</string>
  <key>EnvironmentVariables</key><dict><key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string></dict>
</dict></plist>
