# ops/mac/launchd/ai.tmc.daily_brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.tmc.daily_brief</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/mikicourage/.hermes/hermes-agent/venv/bin/python</string>
    <string>/Users/mikicourage/TradingMapClaw/ops_migration/daily_brief.py</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/Users/mikicourage/TradingMapClaw/ops_migration/logs/daily_brief.out.log</string>
  <key>StandardErrorPath</key><string>/Users/mikicourage/TradingMapClaw/ops_migration/logs/daily_brief.err.log</string>
  <key>EnvironmentVariables</key>
  <dict><key>PATH</key><string>/Users/mikicourage/.local/bin:/usr/local/bin:/usr/bin:/bin</string></dict>
</dict>
</plist>
